# scene_catalog/rootfs/etc/services.d/scene_catalog/run
#!/usr/bin/with-contenv bashio
# Start the Scene Catalog HTTP app

export SCENE_CATALOG_DEFAULT_TARGETS
SCENE_CATALOG_DEFAULT_TARGETS="$(bashio::config 'default_targets')"

exec /usr/bin/python3 /opt/scene_catalog/app.py
